# docker-backup.service
[Unit]
Description=Docker Volume Backup
After=network.target

[Service]
ExecStart=/root/backup_docker.sh
User=root
Type=oneshot
