# docker-backup.timer
[Unit]
Description=Docker Backup

[Timer]
OnCalendar=Mon..Fri 22:00

[Install]
WantedBy=timers.target